# misc/pkcsslotd.service.in
[Unit]
Description=Daemon which manages cryptographic hardware tokens for the openCryptoki package
After=local-fs.target

[Service]
Type=forking
PIDFile=/run/pkcsslotd.pid
ExecStart=@sbindir@/pkcsslotd

[Install]
WantedBy=multi-user.target
